# launch/roles/observer/files/supervisor.conf
[program:observe]
command={{ src }}/observer/./upload.py face
